# conf/systemd-laptop/user/dunst.service.d/override.conf
[Service]
# workaround to make sure the graphical session is ready
ExecStartPre=/bin/sleep 30

[Install]
WantedBy=sway-session.target
