# xdg/systemd-laptop/user/dunst.service.d/override.conf
[Service]
Environment=WAYLAND_DISPLAY=wayland-1
Environment=DISPLAY=:0
# workaround to make sure the graphical session is ready
ExecStartPre=/bin/sleep 30

[Install]
WantedBy=default.target