# xdg/systemd-tubian/user/dev-redis.service
# container-dev-redis.service
# autogenerated by Podman 4.2.0
# Mon Sep 19 14:35:43 CST 2022

# start with podman run -it -p 6379:6379 --name dev-redis docker.io/redis:alpine
# generate this file with podman generate systemd -n dev-redis > ~/.config/systemd/user/dev-redis.service

[Unit]
Description=Podman container-dev-redis.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/user/1000/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start dev-redis
ExecStop=/usr/bin/podman stop -t 10 dev-redis
ExecStopPost=/usr/bin/podman stop -t 10 dev-redis
PIDFile=/run/user/1000/containers/vfs-containers/063ff8b07d63b37fecac201cd9ee32d786d472e7645e349539258cb67ff803da/userdata/conmon.pid
Type=forking

[Install]
WantedBy=default.target
