# xdg/systemd-tubian/user/llama-cpp-server-proxy.socket
[Unit]
Description=Socket activation for llama.cpp server

[Socket]
ListenStream=8080

[Install]
WantedBy=sockets.target
