# xdg/systemd-laptop/user/telegram.service
[Unit]
Description=Telegram

[Service]
ExecStart=/usr/bin/telegram-desktop

[Install]
WantedBy=sway-session.target
