# conf/systemd-laptop/user/fcitx.service
[Unit]
Description=Fcitx input method framework

[Service]
ExecStart=/usr/bin/fcitx5 -r

[Install]
WantedBy=sway-session.target
