# conf/systemd-laptop/user/idle.service
[Unit]
Description=Start swayidle after boot up

[Service]
ExecStart=/home/shou/.config/sway/scripts/idle.sh

[Install]
WantedBy=sway-session.target
