# conf/systemd-laptop/user/kanshi.service
[Unit]
Description=Kanshi monitor configuration auto loader

[Service]
ExecStart=/usr/bin/kanshi

[Install]
WantedBy=sway-session.target
