# conf/systemd-laptop/user/thunderbird.service
[Unit]
Description=Thunderbird

[Service]
ExecStart=/usr/bin/thunderbird-bin

[Install]
WantedBy=sway-session.target
